# single-image-tsdb-grafana/etc/s6-overlay/s6-rc.d/sample-metrics/run
#!/bin/bash
exec /app/venv/bin/python /app/scripts/generate_sample_metrics.py > /app/logs/sample-metrics.log 2>&1
